# docker/rootfs/etc/s6-overlay/s6-rc.d/caddy/run
#!/command/with-contenv sh
/usr/bin/caddy run -c /etc/caddy/Caddyfile
